# .config/autostart/flameshot.desktop
[Desktop Entry]
Type=Application
Name=FlameShot
Exec=flameshot &
#Icon=<full path to icon>
#Comment=<optinal comments>
