# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Name=xBindKeys
Exec=xbindkeys
#Icon=<full path to icon>
#Comment=<optinal comments>
